# systemd/claude-discord-shell.service
[Unit]
Description=Claude Discord Shell
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/davidmarsh/code/claude-discord-shell
EnvironmentFile=/home/davidmarsh/code/claude-discord-shell/.env
ExecStart=/usr/bin/node /home/davidmarsh/code/claude-discord-shell/dist/index.js
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=600
StartLimitBurst=3

[Install]
WantedBy=default.target
