# launchd/com.discoclaw.personal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Reference template for the macOS launchd agent generated by:
    discoclaw install-daemon --service-name personal

  DO NOT load this file directly — the installer bakes your .env variables
  into the plist at install time and writes it to:
    ~/Library/LaunchAgents/com.discoclaw.personal.plist

  To regenerate after a discoclaw update or .env change:
    discoclaw install-daemon --service-name personal
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.discoclaw.personal</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/node</string>
		<string>/path/to/discoclaw/dist/index.js</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/path/to/your/workspace</string>
	<key>EnvironmentVariables</key>
	<dict>
		<!-- Variables from .env are baked in at install time by discoclaw install-daemon -->
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
